# packaging/linux/oakstead.service
[Unit]
Description=Oakstead school records and gradebook server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/oakstead
EnvironmentFile=-/etc/oakstead/oakstead.env
Environment=OAKSTEAD_DATA_DIR=/var/lib/oakstead
Environment=SQLITE_BIN=/usr/bin/sqlite3
Environment=OAKSTEAD_UPDATE_MODE=git
ExecStart=/usr/bin/node /opt/oakstead/server.js
Restart=on-failure
RestartSec=10
User=oakstead
Group=oakstead

[Install]
WantedBy=multi-user.target
